# crontab
@reboot		/usr/sbin/daemon -cf /jenkins/keepalive.sh
*/10 * * * *	/usr/sbin/daemon -cf /jenkins/keepalive.sh
