# crontab
@reboot		/usr/sbin/daemon -cf /usr/local/jenkins/jenkins-slave-scripts/keepalive.sh
*/10 * * * *	/usr/sbin/daemon -cf /usr/local/jenkins/jenkins-slave-scripts/keepalive.sh
