# crontab
@reboot      daemon -cf /home/jenkins/keepalive.sh
*/10 * * * * daemon -cf /home/jenkins/keepalive.sh
