# crontab
@reboot      daemon -cf /jenkins/keepalive.sh
*/10 * * * * daemon -cf /jenkins/keepalive.sh
